# scripts/rag-api.service
[Unit]
Description=FastAPI RAG Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rag
Environment="PATH=/home/ubuntu/rag/venv/bin"
ExecStart=/home/ubuntu/rag/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-api

# 환경 변수 파일 로드
EnvironmentFile=/home/ubuntu/rag/.env

[Install]
WantedBy=multi-user.target
